# packaging/linux/standalone/kymorem-client.service
[Unit]
Description=KyMoRem standalone client
After=graphical-session.target network-online.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=KYMOREM_BIND=0.0.0.0
Environment=KYMOREM_PORT=54865
EnvironmentFile=-%h/.config/kymorem/kymorem.env
ExecStart=@APP_DIR@/run-client.sh
Restart=on-failure
RestartPreventExitStatus=64
RestartSec=2

[Install]
WantedBy=default.target
